# haos-addon/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bash

# Warte bis CouchDB bereit ist
until curl -sf http://127.0.0.1:5984/ > /dev/null 2>&1; do
  echo "[nginx] Warte auf CouchDB..."
  sleep 2
done

# pupils-DB und teacher-User anlegen (idempotent)
COUCHDB_PASSWORD=$(jq -r '.couchdb_password' /data/options.json)
TEACHER_PASSWORD=$(jq -r '.teacher_password' /data/options.json)

curl -sf -u "admin:${COUCHDB_PASSWORD}" -X PUT http://127.0.0.1:5984/pupils || true
curl -sf -u "admin:${COUCHDB_PASSWORD}" -X PUT http://127.0.0.1:5984/_users || true

# teacher-User (CouchDB-User, nicht Admin) — für App-Login
curl -sf -u "admin:${COUCHDB_PASSWORD}" \
  -X PUT http://127.0.0.1:5984/_users/org.couchdb.user:teacher \
  -H "Content-Type: application/json" \
  -d "{\"name\":\"teacher\",\"password\":\"${TEACHER_PASSWORD}\",\"roles\":[],\"type\":\"user\"}" || true

# teacher Lese-/Schreibzugriff auf pupils-DB
curl -sf -u "admin:${COUCHDB_PASSWORD}" \
  -X PUT http://127.0.0.1:5984/pupils/_security \
  -H "Content-Type: application/json" \
  -d '{"admins":{"names":[],"roles":[]},"members":{"names":["teacher"],"roles":[]}}' || true

exec nginx -g "daemon off;"
