# haos-addon/rootfs/etc/services.d/cron/run
#!/command/with-contenv bash
set -a
source /app/backup.env
set +a
exec supercronic /app/crontab
